# systemd/renderis-support-bot.service
[Unit]
Description=Renderis Support Bot
After=network.target

[Service]
Type=simple
User=renderisbot
WorkingDirectory=/opt/renderis-ai-telegram-bot
EnvironmentFile=/etc/renderis-bot.env
ExecStart=/opt/renderis-ai-telegram-bot/.venv/bin/python -m app.support_bot.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
